# .supervisord.conf.erb
[program:proxy-js]
command=node /home/john/chef-solo/proxy.js
autostart=true
autorestart=true
user=root

[program:gnome-screensaver]
command=/usr/bin/gnome-screensaver
autostart=true
autorestart=true
user=john

