# .config/autostart/template.desktop
[Desktop Entry]
Type=Application
Exec=
Hidden=false
X-GNOME-Autostart-enable=true
Name=
Comment=
